# packaging/linux/prerm
#!/bin/sh
set -e

if command -v systemctl >/dev/null 2>&1; then
  TARGET_USER="${SUDO_USER:-$USER}"
  if [ -n "$TARGET_USER" ] && [ "$TARGET_USER" != "root" ]; then
    su - "$TARGET_USER" -c "systemctl --user disable --now print.it.service" 2>/dev/null || true
  fi
fi

exit 0
